# scripts/flask-analyzer.service
[Unit]
Description=HDFS Log Analyzer Flask Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/log-analyzer/scripts
Environment=PATH=/home/YOUR_USERNAME/log-analyzer/scripts/flask-env/bin
ExecStart=/home/YOUR_USERNAME/log-analyzer/scripts/flask-env/bin/python flask-simulator.py
EnvironmentFile=/home/YOUR_USERNAME/log-analyzer/scripts/.env.production
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
